# com.user.ytdlp-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ytdlp-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/yuhaoyong/yt-dlp-extension/server.py</string>
    </array>

    <!-- 环境变量（launchd 默认环境极简，需手动补充） -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/yuhaoyong</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- 登录时自动启动 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 意外退出后自动重启 -->
    <key>KeepAlive</key>
    <true/>

    <!-- 日志输出位置 -->
    <key>StandardOutPath</key>
    <string>/tmp/ytdlp-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ytdlp-server.log</string>
</dict>
</plist>
